# services/upload-scanner.service
[Unit]
Description=OpenClaw Upload Scanner - AI detection on uploaded files
After=network.target

[Service]
Type=simple
User=opc
WorkingDirectory=/home/opc/.openclaw/workspace
ExecStart=/usr/bin/python3 /home/opc/.openclaw/workspace/upload_scanner.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
